# tmpfsmysql.conf
description "Upstart service for starting tmpfs-mysql script"
author "Martin Georgiev"
start on (started networking)
exec bash -c '/home/vagrant/tmpfs-mysql/upstart-service-start.sh'
